# templates/errbot.upstart
#-------------------------------------------------------------------------------
#
# Upstart script for Errbot
#
#-------------------------------------------------------------------------------
# {{ ansible_managed }}
#-------------------------------------------------------------------------------

start on startup
stop on shutdown

respawn

script
	cd {{ errbot_home }}
	exec {{ errbot_bin }}
end script


# vim:ft=upstart:
